# content/etc/supervisord.conf
[supervisord]
pidfile = /run/supervisord.pid
logfile = /var/log/supervisor/supervisord.log
# Set loglevel=debug, only then all logs from child services are printed out
# to container logs (and thus available via `docker logs [container]`
loglevel = debug

nodaemon=true

# These two (unix_http_server, rpcinterface) are needed for supervisorctl to work
[inet_http_server]
port = :9111
username = sv
password = password

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = http://localhost:9111
username = sv
password = password

[include]
files = /etc/supervisor.d/*.conf

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND
redirect_stderr=true
stopasgroup=true

[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch -Des.insecure.allow.root=true --default.path.conf=/etc/elasticsearch/
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true

[program:kibana]
command=/opt/kibana/bin/kibana -c /opt/kibana/config/kibana.yml --log-file /var/log/kibana/kibana.log
autostart=true
autorestart=true
redirect_stderr=true

[program:logstash]
command=/opt/logstash/bin/logstash --config /etc/logstash/conf.d/ --log /var/log/logstash/logstash.log
numprocs=1
autostart=true
autorestart=true
